# templates/init.d.erb
#!/bin/sh
NODE_ENV="production"
NODE_APP='<%=@file_to_run%>'
APP_DIR='<%=@repo_dir%>/<%=@project_name%>';
PID_FILE=/var/run/nodesite_<%=@project_name%>.pid
LOG_FILE=$APP_DIR/app.log
CONFIG_DIR=$APP_DIR/config
#PORT=3000
NODE_EXEC=<%=@node_exec_dir%>/node
NPM_EXEC=<%=@node_exec_dir%>/npm
export HOME="/root"
###############

# REDHAT chkconfig header

# chkconfig: - 58 74
# description: <%=@project_name%> is the script for starting a node app on boot.
### BEGIN INIT INFO
# Provides: node
# Required-Start:    $network $remote_fs $local_fs
# Required-Stop:     $network $remote_fs $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: start and stop nodejs project: <%=@project_name%>
# Description: Node process for $APP_DIR
### END INIT INFO

start_app (){
    if [ -f $PID_FILE ]
    then
        echo "$PID_FILE exists, process is already running or crashed"
        exit 1
    else
        echo "Starting <%=@project_name%> app..."
        #TODO test functionality for port
        #PORT=$PORT 
        cd $APP_DIR
        $NPM_EXEC install
        NODE_ENV=$NODE_ENV NODE_CONFIG_DIR=$CONFIG_DIR $NODE_EXEC $APP_DIR/$NODE_APP <%=@inline_node_params %>  1>$LOG_FILE 2>&1 &
        echo $! > $PID_FILE;
    fi
}

stop_app (){
    if [ ! -f $PID_FILE ]
    then
        echo "$PID_FILE does not exist, process is not running"
        exit 1
    else
        echo "Stopping $APP_DIR/$NODE_APP ..."
        echo "Killing `cat $PID_FILE`"
        kill `cat $PID_FILE`;
        rm -f $PID_FILE;
        echo "<%=@project_name%> stopped"
    fi
}

case "$1" in
    start)
        start_app
    ;;

    stop)
        stop_app
    ;;

    restart)
        stop_app
        start_app
    ;;

    status)
           
     if [ -f $PID_FILE ]; then
        echo "pid file exists"
        PID=`cat $PID_FILE`
        if [  -z "`ps axf | grep ${PID} | grep -v grep`"  ]; then
            echo "Node app stopped but pid file exists. Removing pid: $PID"
            rm -f $PID_FILE
            exit 1
        else
            echo "Node app running with pid $PID"
            exit 0
        fi
    else
        echo "Node app stopped"
        exit 3
    fi
    ;;

    *)
        echo "Usage: $0 {start|stop|restart|status}"
        exit 1
    ;;
esac
